# aws/files/got-quotes-microservice-java/etc/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock                      ; path to your socket file

[supervisord]
logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                          ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                    ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                      ; default user
childlogdir=/var/log/supervisord/              ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock        ; use a unix:// URL  for a unix socket

[program:blogfrontend]
command=/usr/local/bin/ms-run.sh /opt/got-quotes-microservice-java/bin/got-quotes-microservice-java com.schibsted.hackathons.example.gotquotes.AppServer
stderr_logfile = /var/log/got-quotes-microservice-java/got-quotes-microservice-java-stderr.log
stdout_logfile = /var/log/got-quotes-microservice-java/got-quotes-microservice-java-stdout.log
user=root                                      ; default user
autostart=true
autorestart=true
stdout_logfile_maxbytes=100MB
stderr_logfile_maxbytes=100MB
