# aws/files/common/etc/init/datadogLoggly.conf
description "monitoringServices"

start on started elastic-network-interfaces
stop on run level [!2345]
env enabled=1


exec /usr/local/bin/monitor-init.sh
